# crontab.cron
35 12 * * Wed $mydotfiles/cron/brewdeps.sh
36 12 * * Wed $mydotfiles/cron/vimplugins.sh
37 12 * * Wed $mydotfiles/cron/rbenvplugins.sh
