# crontab.cron
35 12 * * Wed $mydotfiles/cron/brewdeps.sh
45 12 * * Wed $mydotfiles/cron/vimplugins.sh
