# crontab.cron
35 12 * * Wed /code/src/utilities/mydotfiles/cron/brewdeps.sh
36 12 * * Wed /code/src/utilities/mydotfiles/cron/vimplugins.sh
37 12 * * Wed /code/src/utilities/mydotfiles/cron/rbenvplugins.sh
