# crontab.cron
35 12 * * Wed /Users/bj/code/src/utilities/dotfiles/cron/brewdeps.sh
36 12 * * Wed /Users/bj/code/src/utilities/dotfiles/cron/vimplugins.sh
